# daemon/systemd/ohmyclawd-daemon.service
[Unit]
Description=ohmyclawd daemon — probes Anthropic for Claude Code utilization
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/ohmyclawd-daemon
Restart=always
RestartSec=5
User=__USER__
Group=__USER__
Environment=OHMYCLAWD_LISTEN=:8787
Environment=OHMYCLAWD_PROBE_INTERVAL=60s
# OHMYCLAWD_CREDS_PATH defaults to ~/.claude/.credentials.json under the target user.
NoNewPrivileges=true
PrivateTmp=false
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=

[Install]
WantedBy=multi-user.target
